# orchestration/launchd/com.ghostbrain.semantic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Weekly semantic-link refresh. Runs Sunday 23:00 (after the metrics
  snapshot at 22:30 + the profile applier at 22:00).

  TEMPLATE — replace placeholders before loading.

  Note: requires the ``[semantic]`` extras to be installed:
    pip install -e ".[semantic]"
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghostbrain.semantic</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/.venv/bin/python</string>
        <string>-m</string>
        <string>ghostbrain.semantic.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/semantic.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/semantic.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_PATH</key>
        <string>__VAULT_PATH__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
