# orchestration/launchd/com.ghostbrain.semantic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Semantic-link refresh, every 15 minutes. Incremental — unchanged notes
  are skipped via mtime/hash, so steady-state runs complete in seconds.
  Worth being aggressive so queries see fresh notes / transcripts.

  TEMPLATE — replace placeholders before loading.

  Note: requires the ``[semantic]`` extras to be installed:
    pip install -e ".[semantic]"
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghostbrain.semantic</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/.venv/bin/python</string>
        <string>-m</string>
        <string>ghostbrain.semantic.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/semantic.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/semantic.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_PATH</key>
        <string>__VAULT_PATH__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
