# orchestration/launchd/com.ghostbrain.jira.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Jira fetcher. Polls every 4 hours per SPEC §10.2.

  TEMPLATE — replace placeholders before loading:
    __REPO_ROOT__   absolute path to your local clone of ghost-brain
    __VAULT_PATH__  absolute path to your vault

  Auth secrets are picked up from .env via python-dotenv on every run.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghostbrain.jira</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/.venv/bin/python</string>
        <string>-m</string>
        <string>ghostbrain.connectors.jira</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>RunAtLoad</key>
    <false/>

    <!-- Every 4 hours. -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/jira.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/jira.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_PATH</key>
        <string>__VAULT_PATH__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
